# packaging/deb/postinst
#!/bin/bash

set -e

# Create tproxy user and group if they don't exist
if ! getent group tproxy >/dev/null; then
    groupadd --system tproxy
fi

if ! getent passwd tproxy >/dev/null; then
    useradd --system --no-create-home --shell /bin/false --gid tproxy tproxy
fi

# Create log directory
mkdir -p /var/log/tproxy
chown tproxy:tproxy /var/log/tproxy
chmod 755 /var/log/tproxy

# Set permissions on config directory
chown -R tproxy:tproxy /etc/tproxy
chmod 755 /etc/tproxy
chmod 644 /etc/tproxy/*.yaml

# Reload systemd to pick up new service
systemctl daemon-reload

# Enable service to start on boot (but don't start it automatically)
systemctl enable tproxy.service

echo "tproxy service installed. To start the service, run:"
echo "  systemctl start tproxy"
echo ""
echo "To check the status, run:"
echo "  systemctl status tproxy"